# .config/systemd/user/change_wallpaper.service
[Unit]
Description=Change wallpaper using SWWW
After=swww.service

[Service]
ExecStart=/bin/bash /home/esalexreyes/.config/swww/swww.sh

[Install]
WantedBy=default.target
